# api-gateway/backend/scripts/news-fetcher.service
[Unit]
Description=News Fetcher Proxy (Flask)
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root
ExecStart=/usr/bin/python3 /root/news_fetcher.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
